# config/infrastructure/php/supervisord.conf
[supervisord]
nodaemon=true

[program:blackfire-agent]
command=/usr/bin/blackfire-agent --socket=unix:///var/run/blackfire/agent.sock

[program:php-fpm]
command=docker-php-entrypoint php-fpm
